# etc/cron.d/clean-system
# Clean System cron job
15 04 * * * root /usr/local/bin/clean-system --raw >/dev/null 2>&1
