# etc/cron.d/update-ip
# update the IP address
*/5 * * * * root /usr/local/share/CasjaysDev/scripts/bin/update-ip --raw >/dev/null 2>&1
