# cron.d/fun-stuff
#Lets have fortune and cowsay create an issue file 30 minutes
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games
*/15,45 * * * * root fortune | cowsay > /etc/motd.net && printf "\n\n" >> /etc/motd.net ; cat /etc/casjaysdev/messages/legal.txt /etc/motd.net > /etc/issue.net
